# source/cabal-data/cabal_scripts/cabal_services/cabal_server3.conf
[program:DBAgent_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/DBAgent_03'
user=root
autostart=false
autorestart=true
startsecs=3

[program:CashDBAgent_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/CashDBAgent_03'
user=root
autostart=false
autorestart=true
startsecs=3

[program:PartySvr_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/PartySvr_03'
user=root
autostart=false
autorestart=true
startsecs=3

[program:ChatNode_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/ChatNode_03'
user=root
autostart=false
autorestart=true
startsecs=3

[program:AgentShop_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/AgentShop_03'
user=root
autostart=false
autorestart=true
startsecs=3




[program:EventDBAgent_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/EventDBAgent_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:EventMgrSvr_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/EventMgrSvr_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3





[program:WorldSvr_03_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_01'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_02'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_03'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_04]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_04'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_05]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_05'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_06]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_06'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_07]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_07'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_08]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_08'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_09]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_09'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_10]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_10'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_11]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_11'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_12]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_12'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_13]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_13'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_14]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_14'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_15]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_15'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_16]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_16'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_17]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_17'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_18]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_18'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_19]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_19'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_20]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_20'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_21]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_21'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_22]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_22'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_23]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_23'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_24]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_24'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_25]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_25'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_26]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_26'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_27]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_27'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_28]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_28'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_29]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_29'
user=root
autostart=false
autorestart=true
startsecs=15

[program:WorldSvr_03_30]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_03_30'
user=root
autostart=false
autorestart=true
startsecs=15
