# source/cabal-data/cabal_scripts/cabal_services/cabal_server1.conf
[program:DBAgent_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/DBAgent_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:CashDBAgent_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/CashDBAgent_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:PartySvr_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/PartySvr_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:ChatNode_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/ChatNode_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:AgentShop_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/AgentShop_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3




[program:EventDBAgent_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/EventDBAgent_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:EventMgrSvr_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/EventMgrSvr_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3




[program:WorldSvr_01_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_03 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_04]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_04 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_05]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_05 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_06]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_06 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_07]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_07 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_08]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_08 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_09]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_09 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_10]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_10 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_11]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_11 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_12]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_12 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_13]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_13 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_14]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_14 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_15]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_15 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_16]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_16 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_17]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_17 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_18]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_18 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_19]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_19 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_20]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_20 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_21]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_21 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_22]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_22 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_23]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_23 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_24]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_24 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_25]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_25 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_26]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_26 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_27]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_27 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_28]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_28 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_29]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_29 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_01_30]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_01_30 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15
