# source/cabal-data/cabal_scripts/cabal_services/cabal_server2.conf
[program:DBAgent_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/DBAgent_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:CashDBAgent_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/CashDBAgent_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:PartySvr_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/PartySvr_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:ChatNode_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/ChatNode_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:AgentShop_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/AgentShop_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3




[program:EventDBAgent_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/EventDBAgent_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3

[program:EventMgrSvr_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/EventMgrSvr_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=3





[program:WorldSvr_02_01]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_01 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_02]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_02 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_03]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_03 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_04]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_04 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_05]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_05 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_06]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_06 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_07]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_07 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_08]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_08 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_09]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_09 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_10]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_10 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_11]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_11 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_12]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_12 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_13]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_13 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_14]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_14 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_15]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_15 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_16]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_16 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_17]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_17 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_18]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_18 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_19]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_19 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_20]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_20 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_21]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_21 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_22]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_22 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_23]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_23 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_24]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_24 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_25]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_25 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_26]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_26 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_27]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_27 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_28]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_28 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_29]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_29 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15

[program:WorldSvr_02_30]
directory=/etc/cabal_etc/core/
command=/bin/bash -c 'exec /usr/bin/WorldSvr_02_30 -l -c'
user=root
autostart=false
autorestart=unexpected
startsecs=15
